# target/linux/ramips/mt7621/base-files/etc/uci-defaults/05_fix-compat-version
#
# Copyright (C) 2020 OpenWrt.org
#

. /lib/functions.sh
board=$(board_name)

case "$board" in
	d-team,newifi-d2|\
	gehua,ghl-r-001|\
	netgear,r6220|\
	netgear,wndr3700-v5|\
	xiaomi,mir3g-nor|\
	xiaomi,mir3g|\
	xiaomi,mi-router-ac2100|\
	xiaomi,redmi-router-ac2100|\
	xwrt,nxc200p|\
	xwrt,puppies|\
	xwrt,t-cpe1201k-v01|\
	xwrt,t-cpe1200k-v01|\
	xwrt,t-cpe1202kd-v01|\
	xwrt,x-sdwan-1200)
		exit 0
	;;
esac

uci set system.@system[0].compat_version="1.1"
uci commit system

exit 0
